# roles/security/crowdsec/templates/upgrade.service
[Unit]
Description=CrowdSec Hub Upgrade Service

[Service]
Type=oneshot
ExecStart=/usr/bin/cscli hub upgrade